# etc/download.timer
[Unit]
Description=netfilter download every 10 minutes

[Timer]
OnActiveSec=0
OnCalendar=*:0/10

[Install]
WantedBy=multi-user.target

